# scripts/supervisor/conf.d/zenqa.ini
[program:zenqa]
command=pipenv run gunicorn zenofqa.wsgi -w 2 -k gthread -b 127.0.0.1:8000
directory=/home/devin/apps/zenqa
autostart=true
autorestart=unexpected
user=devin
stdout_logfile=/home/devin/apps/zenqa/scripts/supervisor/var/log/zenqa-stdout.log
stderr_logfile=/home/devin/apps/zenqa/scripts/supervisor/var/log/zenqa-stderr.log
